# prometheus.service
[Unit]
Description=Prometheus

[Service]
ExecStart=/opt/prometheus/prometheus.sh

[Install]
WantedBy=multi-user.target
