# prometheus.service
[Unit]
Description=Prometheus

[Service]
ExecStart=/opt/prometheus/prometheus

[Install]
WantedBy=multi-user.target
